# templates/launchd-folder-watch.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
Agent Bootstrap — Folder Watch launchd
파일: ~/Library/LaunchAgents/com.{USERNAME}.folder-watch.plist
WATCH_FOLDER 에 파일 변화 감지 시 run.sh 실행.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.{USERNAME}.folder-watch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{HOME}/.claude/scripts/folder-watch/run.sh</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>{WATCH_FOLDER}</string>
    </array>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>{HOME}/.claude/scripts/folder-watch/logs/launchd.out</string>

    <key>StandardErrorPath</key>
    <string>{HOME}/.claude/scripts/folder-watch/logs/launchd.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
