# templates/systemd-folder-watch.path.tmpl
[Unit]
Description=Agent Bootstrap — Watch {WATCH_FOLDER}

[Path]
PathChanged={WATCH_FOLDER}
TriggerLimitIntervalSec=5

[Install]
WantedBy=default.target
