# templates/systemd-slack-jipsa.service.tmpl
# Agent Bootstrap — Slack Jipsa daemon (systemd user unit)
#
# 파일명:
#   글로벌:    ~/.config/systemd/user/slack-jipsa.service
#   프로젝트:  ~/.config/systemd/user/slack-jipsa-{PROJECT_ID}.service
#
# 플레이스홀더:
#   {INSTANCE}  글로벌이면 "slack-jipsa", 프로젝트 모드면 "slack-jipsa-{PROJECT_ID}"
[Unit]
Description=Agent Bootstrap — Slack Jipsa daemon ({INSTANCE})
After=default.target

[Service]
Type=simple
Environment=SLACK_JIPSA_INSTANCE={INSTANCE}
EnvironmentFile=%h/.claude/secrets/{INSTANCE}.env
ExecStart=/usr/bin/python3 %h/.claude/scripts/{INSTANCE}/daemon.py
Restart=always
RestartSec=5
StandardOutput=append:%h/.claude/scripts/{INSTANCE}/logs/systemd.out
StandardError=append:%h/.claude/scripts/{INSTANCE}/logs/systemd.err

[Install]
WantedBy=default.target
